# docker/cron/crontab
* * * * * /usr/local/bin/php /var/www/bin/magento cron:run >> /var/www/var/log/run.cron.log
